# deploy/ansible/roles/portman/templates/portman-gc.timer.j2
[Unit]
Description=Run portman gc periodically

[Timer]
OnCalendar={{ portman_gc_interval | default('hourly') }}
Persistent=true

[Install]
WantedBy=timers.target
