# deploy/ansible/roles/portman/templates/portman-gc.service.j2
[Unit]
Description=portman garbage collection

[Service]
Type=oneshot
ExecStart=/usr/local/bin/portman gc
